# aviation_feeder/rootfs/etc/s6-overlay/s6-rc.d/sdrmap-stunnel/run
#!/command/with-contenv bash
# stunnel TLS tunnel for sdrmap MLAT (apt-installed stunnel at /usr/bin/stunnel):
# listens on 127.0.0.1:3333 and forwards, TLS-wrapped, to mlat.feed.sdrmap.org.
# sdrmap-mlat's mlat-client connects to :3333. Only run when sdrmap MLAT is
# actually in play (enabled + credentials + a station position); idle otherwise.
. /usr/local/bin/feeder-gate
gate "[sdrmap-stunnel]" ENABLE_SDRMAP SMUSERNAME SMPASSWORD LAT LONG ALT
exec /usr/local/bin/s6wrap --quiet --timestamps --prepend=sdrmap-stunnel --args \
  /usr/bin/stunnel /etc/stunnel/mlat.conf
