# aviation_feeder/rootfs/etc/s6-overlay/s6-rc.d/fr24feed/run
#!/command/with-contenv sh
# shellcheck shell=sh
# Gate the FlightRadar24 1090 feeder: needs enable_fr24 + a sharing key.
# Without the key, 01-fr24feed skips writing fr24feed.ini, so the binary would
# otherwise restart-loop.
. /usr/local/bin/feeder-gate
gate "[fr24feed]" ENABLE_FR24 FR24KEY
# fr24feed's binary demands TZ=GMT, but the CONTAINER stays on the user's tz. We do
# NOT force GMT on this whole run script: that also put s6wrap (the --timestamps
# logger) on GMT, so fr24's log lines printed 4h off from every sibling service. The
# GMT requirement is pinned to the binary ALONE by the /usr/local/bin/fr24feed
# wrapper (see Dockerfile), so this script -- and its s6wrap logger -- run on the
# container tz and stamp fr24's logs in local time.
exec /etc/s6-overlay/scripts/fr24feed
