# aviation_feeder/rootfs/etc/s6-overlay/s6-rc.d/dump978/run
#!/command/with-contenv sh
# shellcheck shell=sh
# Gate the 978 UAT decoder: run it whenever 978 is decoded locally -- rtlsdr mode
# with UAT enabled, or uat-only mode (00-haos-options forces ENABLE_UAT on there).
# In remote mode there is no local SDR, and UAT arrives via the remote Beast.
if [ "${ENABLE_UAT:-false}" != "true" ] || [ "${RECEIVER_MODE:-rtlsdr}" = "remote" ]; then
  echo "[dump978] UAT decode off (enable_uat off or remote mode); idling"
  exec sleep infinity
fi
exec /etc/s6-overlay/scripts/dump978
