# aviation_feeder/rootfs/etc/s6-overlay/s6-rc.d/pfclient/run
#!/command/with-contenv sh
# Gate PlaneFinder: needs enable_planefinder + a sharecode. Idle otherwise.
. /usr/local/bin/feeder-gate
gate "[pfclient]" ENABLE_PLANEFINDER SHARECODE
exec /etc/s6-overlay/scripts/pfclient
