# aviation_feeder/rootfs/etc/s6-overlay/s6-rc.d/fr24uat-feed/run
#!/command/with-contenv sh
# shellcheck shell=sh
# Gate the FlightRadar24 UAT (978) feeder: needs enable_fr24 plus a separate
# UAT sharing key (fr24_uat_key). Idle otherwise so it never error-loops.
. /usr/local/bin/feeder-gate
gate "[fr24uat-feed]" ENABLE_FR24 FR24KEY_UAT
# Same fr24feed binary as the 1090 feeder, same GMT requirement -- pinned to the
# binary by the /usr/local/bin/fr24feed wrapper (see Dockerfile), NOT forced here,
# so this script's s6wrap logger stays on the container tz and stamps UAT logs in
# local time rather than GMT.
exec /etc/s6-overlay/scripts/fr24uat-feed
