# aviation_feeder/rootfs/etc/s6-overlay/s6-rc.d/piaware/run
#!/command/with-contenv sh
# Gate the piaware feeder behind the enable_piaware add-on option. When
# disabled the longrun stays "up" but idle so it never connects to
# FlightAware. Overrides the thin run shim copied from docker-piaware.
. /usr/local/bin/feeder-gate
gate "[piaware]" ENABLE_PIAWARE
exec /etc/s6-overlay/scripts/piaware
